# crontab
@reboot autossh -M 0 -f -T -N my-tunnel
